# deploy/rendering-deploy.service
[Unit]
Description=Rendering Automated Deployment Job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=RENDERING_DEPLOY_ENV=/etc/rendering/rendering-deploy.env
WorkingDirectory=/srv/rendering
# git reset/clean 可能覆盖脚本的可执行权限，每次运行前确保权限正确
ExecStartPre=/bin/chmod +x /srv/rendering/deploy/rendering-deploy.sh
ExecStart=/srv/rendering/deploy/rendering-deploy.sh
StandardOutput=journal
StandardError=journal